# systemd/xray-auto-update.service
[Unit]
Description=Xray Auto-Update Service
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/xray-auto-update.sh
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target